# ubuntu_srv/mybootcmd.service
[Unit]
Description=myboot service

[Service]
ExecStart=/etc/init.d/myscript.sh

[Install]
WantedBy=multi-user.target
